# rootfs/etc/services.d/keeper/run
#!/usr/bin/with-contenv bash
# TODO: actually start the keeper here

if [ -z "$RPC_HOST" ]; then
    echo "RPC_HOST not set"
    exit 1
fi

# TODO: wait for consul

# TODO: read addresses out of consul
ETH_FROM=0x0101010101010101010101010101010101010101

# TODO: make sure transactions aren't stuck
# makerdao-helper plunger --source parity_txqueue --wait "$ETH_FROM"

# TODO: how should we pick what keeper to run and with what flags? each should maybe have its own start script
tail -f /dev/null
