# rootfs/etc/services.d/keeper/run
#!/usr/bin/with-contenv bash

# TODO: actually start the keeper here

tail -f /dev/null
